# daemon/com.stockpulse.claudeheartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockpulse.claudeheartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/rajneeshmishra/Downloads/stock-pulse/daemon/claude_heartbeat.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rajneeshmishra/Downloads/stock-pulse</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/rajneeshmishra/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/rajneeshmishra</string>
    </dict>

    <!-- Heartbeat fires every hour; it's a one-shot per firing (not KeepAlive) -->
    <key>RunAtLoad</key>
    <false/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/rajneeshmishra/Downloads/stock-pulse/logs/claude_heartbeat.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rajneeshmishra/Downloads/stock-pulse/logs/claude_heartbeat.err.log</string>
</dict>
</plist>
